# etc/rc.d/rc.local
#! /bin/sh

# $Id$

# Put something useful here
