# trunk/etc/rc.d/rc.local
#! /bin/sh

# $Id$

# Send a system startup email
/usr/local/sbin/startup.sh
